# ansible/roles/tailscale_mgmt/templates/tailscaled-mgmt.service.j2
[Unit]
Description=Tailscale management daemon (second instance, persistent)
# This unit runs a DEDICATED tailscaled for persistent management access.
# The DEFAULT tailscaled.service instance remains free for ephemeral per-job
# joins (tag:ci-deploy via tailscale/github-action) and must NOT be touched
# by this unit.
#
# Key differences from the default tailscaled.service:
#   --statedir  → separate state directory ({{ tailscale_mgmt_statedir }})
#   --socket    → separate Unix socket ({{ tailscale_mgmt_socket }})
#   --tun       → separate TUN interface ({{ tailscale_mgmt_tun }})
#   --port      → separate UDP port ({{ tailscale_mgmt_port }}, 0 = OS-assigned)
#
# Reference: Tailscale docs on running multiple tailscaled instances
# (separate --state/--socket/--tun/--port per instance).
# tailscaled.service (default instance) is started by the tailscale role
# with authkey="" — installed and running but unregistered. Ordering after
# it ensures the TUN subsystem is ready before the mgmt instance starts.
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=notify
ExecStartPre=/usr/sbin/tailscaled --cleanup --socket={{ tailscale_mgmt_socket }}
ExecStart=/usr/sbin/tailscaled \
  --statedir={{ tailscale_mgmt_statedir }} \
  --socket={{ tailscale_mgmt_socket }} \
  --tun={{ tailscale_mgmt_tun }} \
  --port={{ tailscale_mgmt_port }}
ExecStopPost=/usr/sbin/tailscaled --cleanup --socket={{ tailscale_mgmt_socket }}

Restart=on-failure
RestartSec=5

# Runtime socket directory
# RuntimeDirectory creates /run/tailscale-mgmt/ at service start and removes
# it at stop. Ansible does NOT pre-create this directory — systemd owns it.
RuntimeDirectory=tailscale-mgmt
RuntimeDirectoryMode=0750

# Persistent state directory — REQUIRED with ProtectSystem=strict, which
# otherwise mounts /var read-only and would prevent tailscaled from writing
# its node key/prefs to --statedir (losing identity on every restart).
# StateDirectory= creates /var/lib/tailscale-mgmt and marks it read-write for
# this unit. systemd owns it; Ansible does NOT pre-create it.
StateDirectory=tailscale-mgmt
StateDirectoryMode=0750

# Hardening. NOTE: the upstream tailscaled.service does NOT use
# ProtectSystem=strict; we add it here for defence in depth. tailscaled needs
# RW access to its state dir (covered by StateDirectory above) and manages
# kernel netfilter/routes (not blocked by ProtectSystem). Re-verify with a
# real `tailscale up` if upstream hardening changes.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=true
NoNewPrivileges=false
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
