# ansible/roles/tailscale_mgmt/templates/tailscaled-mgmt.service.j2
[Unit]
Description=Tailscale management daemon (second instance, persistent)
# This unit runs a DEDICATED tailscaled for persistent management access.
# The DEFAULT tailscaled.service instance remains free for ephemeral per-job
# joins (tag:ci-deploy via tailscale/github-action) and must NOT be touched
# by this unit.
#
# Key differences from the default tailscaled.service:
#   --statedir  → separate state directory ({{ tailscale_mgmt_statedir }})
#   --socket    → separate Unix socket ({{ tailscale_mgmt_socket }})
#   --tun       → separate TUN interface ({{ tailscale_mgmt_tun }})
#   --port      → separate UDP port ({{ tailscale_mgmt_port }}, 0 = OS-assigned)
#
# Reference: Tailscale docs on running multiple tailscaled instances
# (separate --state/--socket/--tun/--port per instance).
# tailscaled.service (default instance) is started by the tailscale role
# with authkey="" — installed and running but unregistered. Ordering after
# it ensures the TUN subsystem is ready before the mgmt instance starts.
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=notify
ExecStartPre=/usr/sbin/tailscaled --cleanup --socket={{ tailscale_mgmt_socket }}
ExecStart=/usr/sbin/tailscaled \
  --statedir={{ tailscale_mgmt_statedir }} \
  --socket={{ tailscale_mgmt_socket }} \
  --tun={{ tailscale_mgmt_tun }} \
  --port={{ tailscale_mgmt_port }}
ExecStopPost=/usr/sbin/tailscaled --cleanup --socket={{ tailscale_mgmt_socket }}

Restart=on-failure
RestartSec=5

# Runtime socket directory
# RuntimeDirectory creates /run/tailscale-mgmt/ at service start and removes
# it at stop. Ansible does NOT pre-create this directory — systemd owns it.
RuntimeDirectory=tailscale-mgmt
RuntimeDirectoryMode=0750

# Persistent state directory: systemd creates /var/lib/tailscale-mgmt and owns
# it (Ansible does NOT pre-create it). tailscaled writes its node key/prefs here.
StateDirectory=tailscale-mgmt
StateDirectoryMode=0750

# NO sandbox/capability hardening here — DELIBERATE.
# This tailscaled runs `--ssh`, i.e. it is a Tailscale-SSH SERVER: it spawns
# login sessions for the target user, which requires CAP_SETUID/CAP_SETGID
# (initgroups/setuid), write access to /run (motd) and /home, and the ability
# to gain privileges. Restricting these — e.g. CapabilityBoundingSet without
# SETUID/SETGID, ProtectSystem=strict, ProtectHome=yes — breaks every SSH
# session with "initgroups: Operation not permitted" / read-only /run.
# The upstream tailscaled.service runs unrestricted for the same reason; match it.

[Install]
WantedBy=multi-user.target
